# src/pkg_recovery/package-recovery.service
[Unit]
Description=Installation Recovery Service
Requires=dbus.service

[Service]
ExecStart=/usr/bin/pkg_recovery
