# src/pkgdir_tool/pkgdir-tool.service
[Unit]
Description=User Directory Creator

[Service]
ExecStart=/usr/bin/pkgdir-tool
