# src/pkgdir_tool/pkgdir-tool.service
[Unit]
Description=User Directory Creator

[Service]
User=app_fw
Group=app_fw
SmackProcessLabel=System
ExecStart=/usr/bin/pkgdir-tool
